# root/etc/services.d/01-rclone-mount-nocache/finish
#!/usr/bin/with-contenv /bin/sh 

fusermount -u -z /mnt/godrive